# hookd.service
[Unit]
Description=hookd — webhook and schedule dispatcher
Documentation=file:///home/rocky/hookd/DAEMON.md
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=rocky
Group=rocky
WorkingDirectory=/home/rocky/hookd
ExecStart=/usr/bin/python3 /home/rocky/hookd/hookd.py \
    --config /home/rocky/hookd/config.yml
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hookd
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/home/rocky/hookd
AmbientCapabilities=CAP_SETUID CAP_SETGID
CapabilityBoundingSet=CAP_SETUID CAP_SETGID

[Install]
WantedBy=multi-user.target
